# etc-systemd-system/ncbkp.service
[Unit]
Description=Nextcloud backup

[Service]
User=http
Group=http
ExecStart=/bin/bash /path/to/ncbkp.sh
